# pia_service/templates/pia-pf-renew.service.jinja
[Unit]
Description=Renew port forward for PIA VPN

[Service]
Type=oneshot
User={{ user }}
ExecStart={{ pia_service }} renew-port
