# pia_service/templates/pia-pf-renew.timer
[Unit]
Description=Auto-renew port forward for PIA VPN

[Timer]
OnActiveSec=15m
OnUnitActiveSec=15m
