# pia_service/templates/pia-pf-renew.service.jinja
[Unit]
Description=Renew port forward for PIA VPN

[Service]
Type=oneshot
ExecStart={{ python_interpreter }} {{ package_dir }}/cli.py renew-port
